# systemd/jarvis-tts.service
[Unit]
Description=Джарвис — «Голос» (TTS: Silero, кэш WAV)
# Аудио-сервис: нужен пользовательский PipeWire
Wants=pipewire.service
After=pipewire.service network.target
# Переживаем крэш-бурст, но жёсткий цикл падений ограничиваем (виден, не маскируется).
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
# ПРИМЕР. Реальные юниты с верными путями ставит `jarvis start` (в ~/.config/systemd/user/).
# Здесь %h/jarvis/.venv — лишь образец, если проект лежит в ~/jarvis.
# Рабочая директория — корень проекта (страховка для относительных путей).
WorkingDirectory=%h/jarvis
ExecStart=%h/jarvis/.venv/bin/jarvis-tts
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
