# systemd/jarvis-tts.service
[Unit]
Description=Джарвис — «Голос» (TTS: Piper)
# Аудио-сервис: нужен пользовательский PipeWire
Wants=pipewire.service
After=pipewire.service network.target

[Service]
Type=simple
# ПРИМЕР. Реальные юниты с верными путями ставит `jarvis start` (в ~/.config/systemd/user/).
# Здесь %h/jarvis/.venv — лишь образец, если проект лежит в ~/jarvis.
ExecStart=%h/jarvis/.venv/bin/jarvis-tts
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
