# systemd/jarvis-os-agent.service
[Unit]
Description=Джарвис — «Руки» (OS-агент)
After=network.target

[Service]
Type=simple
# ПРИМЕР. Реальные юниты с верными путями ставит `jarvis start` (в ~/.config/systemd/user/).
# Здесь %h/jarvis/.venv — лишь образец, если проект лежит в ~/jarvis.
ExecStart=%h/jarvis/.venv/bin/jarvis-os-agent
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
